# deploy/schoolbell-daemon.service
# VOORBEELD — install.sh genereert de echte unit met jouw paden.
# Handmatig gebruiken? Vervang <gebruiker> en /pad/naar/schoolbell.
# Zie deploy/README.md.
[Unit]
Description=Schoolbell Daemon
After=network.target schoolbell-web.service
Wants=network.target schoolbell-web.service

[Service]
Type=simple
User=<gebruiker>
WorkingDirectory=/pad/naar/schoolbell

EnvironmentFile=/etc/schoolbell/daemon.env

Environment="SCHOOLBELL_CONFIG=/etc/schoolbell/config.json"
Environment="API_BASE=http://127.0.0.1:5000"
# SCHOOLBELL_WEB_USER en SCHOOLBELL_WEB_PASS komen uit daemon.env

ExecStart=/pad/naar/schoolbell/venv/bin/python -u /pad/naar/schoolbell/schoolbelldaemon.py

# `systemctl reload schoolbell-daemon` stuurt SIGHUP: de daemon
# herlaadt dan config.json zonder herstart (en zonder een lopende
# bel te onderbreken). Sinds de mtime-watch pikt hij wijzigingen
# via de webinterface al vanzelf op; dit is de expliciete knop voor
# handmatige edits aan /etc/schoolbell/config.json.
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=2

# Netjes stoppen
TimeoutStopSec=15
KillSignal=SIGTERM
KillMode=mixed

[Install]
WantedBy=multi-user.target
